# init.d/s6/zapret/up
#!/bin/execlineb -P
exec /opt/zapret/init.d/sysv/zapret start
