# init.d/runit/zapret/finish
#!/bin/sh
/opt/zapret/init.d/sysv/zapret stop
