# init.d/runit/zapret/run
#!/bin/sh
/opt/zapret/init.d/sysv/zapret start
exec chpst -b zapret sleep infinity
